# 99-play-disc-to-chrome-cast.rules
SUBSYSTEM=="block", KERNEL=="sr0", ENV{ID_CDROM_MEDIA_CD}=="1", ACTION=="remove", RUN+=CHANGE_THIS_DIR/stop.sh"
